# supervisord.conf
; supervisor config file

[unix_http_server]
file=/home/microk8s/log/supervisor.sock   ; (the path to the socket file)
chmod=0700                       ; sockef file mode (default 0700)

[supervisord]
logfile=/home/microk8s/log/supervisord.log ; (main log file;default $CWD/supervisord.log)
pidfile=/home/microk8s/log/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
childlogdir=/home/microk8s/log/supervisor            ; ('AUTO' child log dir, default $TEMP)

; the below section must remain in the config file for RPC
; (supervisorctl/web interface) to work, additional interfaces may be
; added by defining them in separate rpcinterface: sections
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///home/microk8s/log/supervisor.sock ; use a unix:// URL  for a unix socket

; The [include] section can just contain the "files" setting.  This
; setting can list multiple files (separated by whitespace or
; newlines).  It can also contain wildcards.  The filenames are
; interpreted as relative to this file.  Included files *cannot*
; include files themselves.

#[include]
#files = /etc/supervisor/conf.d/*.conf

###########################################################
# dblog
[group:celery]
programs=dblog-worker, dblog-beat

[program:dblog-worker]
command                 = /home/microk8s/pyenv/versions/dblog/bin/celery  -A  dblog  worker -P gevent -c 2  --loglevel=info
directory               = /home/microk8s/code/dblog
autostart               = false
autorestart             = false
redirect_stderr         = true
stdout_logfile          = /home/microk8s/log/dblog-worker.log
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups  = 1

[program:dblog-beat]
command                 = /home/microk8s/pyenv/versions/dblog/bin/celery  -A  dblog  beat  --loglevel=info
directory               = /home/microk8s/code/dblog
autostart               = false
autorestart             = false
redirect_stderr         = true
stdout_logfile          = /home/microk8s/log/dblog-beat.log
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups  = 1

[program:dblog-service]
command                 = /home/microk8s/pyenv/versions/dblog/bin/python /home/microk8s/code/dblog/script/service.py
directory               = /home/microk8s/code/dblog
autostart               = false
autorestart             = false
redirect_stderr         = true
stdout_logfile          = /home/microk8s/log/dblog-service.log
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups  = 1

[program:dblog-web]
command                 = /home/microk8s/pyenv/versions/dblog/bin/gunicorn -c g.py dblog.wsgi:application
directory               = /home/microk8s/code/dblog
autostart               = true
autorestart             = true
redirect_stderr         = true
stdout_logfile          = /home/microk8s/log/dblog-access.log
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups  = 1


